# pi-loc.conf
# This is the supervisor config file.  It goes in /etc/supervisor/conf.d
[program:pi-loc]
command=/usr/bin/python /opt/pi-loc/main.py
autorestart=true
directory=/opt/pi-loc